# extra/coturn/turnserver.service
[Unit]
Description=STUN and TURN relay server for VoIP and WebRTC
After=network-online.target
After=remote-fs.target
Wants=network-online.target

[Service]
User=turnserver
Type=notify
ExecStart=/usr/bin/turnserver -c /etc/turnserver/turnserver.conf -v
LogsDirectory=turnserver

[Install]
WantedBy=multi-user.target
