# extra/linuxconsole/inputattach.service
[Unit]
Description=Attach serial input devices to kernel input subsystem

[Service]
Type=forking
ExecStart=/usr/bin/inputattachctl start
ExecStop=/usr/bin/inputattachctl stop

[Install]
WantedBy=multi-user.target
