# extra/glances/glances-web.service
[Unit]
Description=Glances Server
After=network.target

[Service]
EnvironmentFile=/etc/conf.d/glances
ExecStart=/usr/bin/glances -w $GLANCES_ARGS
Restart=always
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
