# extra/sssd/sssd-perms.service
[Service]
ExecStartPre=+-/usr/bin/chmod -R u+rw,o-rwx /etc/sssd
ExecStartPre=+-/usr/bin/chgrp sssd /usr/lib/sssd/sssd/krb5_child
ExecStartPre=+-/usr/bin/setcap cap_dac_read_search,cap_setgid,cap_setuid=p /usr/lib/sssd/sssd/krb5_child
ExecStartPre=+-/usr/bin/chgrp sssd /usr/lib/sssd/sssd/ldap_child
ExecStartPre=+-/usr/bin/setcap cap_dac_read_search=p /usr/lib/sssd/sssd/ldap_child
ExecStartPre=+-/usr/bin/chgrp sssd /usr/lib/sssd/sssd/sssd_pam
ExecStartPre=+-/usr/bin/setcap cap_dac_read_search=p /usr/lib/sssd/sssd/sssd_pam
