# extra/tlsrpt-reporter/tlsrpt-collectd.service
[Unit]
Description=tlsrpt-collectd agent to collect TLSRPT reports
After=network.target

[Service]
ExecStart=/usr/bin/tlsrpt-collectd
PrivateDevices=true
PrivateTmp=true
ProtectSystem=true
Restart=always
User=tlsrpt
Group=postfix
LogsDirectory=tlsrpt
StateDirectory=tlsrpt
RuntimeDirectory=tlsrpt
ConfigurationDirectory=tlsrpt

[Install]
WantedBy=multi-user.target
