# extra/tlsrpt-reporter/tlsrpt-reportd.service
[Unit]
Description=tlsrpt-reportd agent to create TLSRPT reports
After=network.target

[Service]
ExecStart=/usr/bin/tlsrpt-reportd
PrivateDevices=true
PrivateTmp=true
ProtectSystem=true
Restart=always
User=tlsrpt
LogsDirectory=tlsrpt
StateDirectory=tlsrpt
RuntimeDirectory=tlsrpt
ConfigurationDirectory=tlsrpt

[Install]
WantedBy=multi-user.target
