# extra/garage/systemd.service
[Unit]
Description=Garage Data Store
Documentation=https://garagehq.deuxfleurs.fr/documentation/quick-start/
After=network-online.target
Wants=network-online.target

[Service]
User=garage
Group=garage
NoNewPrivileges=true
Environment='RUST_LOG=garage=info' 'RUST_BACKTRACE=1'
ExecStart=/usr/bin/garage -c /etc/garage.toml server
Restart=on-failure
RestartSec=20s
StartLimitInterval=30min
StartLimitBurst=30

[Install]
WantedBy=multi-user.target
