# extra/zabbix/zabbix-agent2.service
[Unit]
Description=Zabbix Agent 2
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/zabbix_agent2
User=zabbix-agent

[Install]
WantedBy=multi-user.target
