# extra/zabbix/zabbix-web-service.service
[Unit]
Description=Zabbix Web Service
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/zabbix_web_service -f
User=zabbix

[Install]
WantedBy=multi-user.target
