# extra/zabbix/zabbix-proxy-pgsql.service
[Unit]
Description=Zabbix proxy with PostgreSQL
After=network-online.target postgresql.service
Conflicts=zabbix-proxy-sqlite.service zabbix-proxy-mysql.service

[Service]
Type=simple
ExecStart=/usr/bin/zabbix_proxy_postgresql -f
User=zabbix
PrivateDevices=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
CapabilityBoundingSet=cap_net_raw
AmbientCapabilities=cap_net_raw

[Install]
WantedBy=multi-user.target
